# asus-xonar@.service
[Unit]
Description=ASUS Xonar daemon
BindTo=%i.device
After=%i.device

[Service]
Type=forking
TimeoutStartSec=5
ExecStart=/usr/local/bin/xonard /%I
